# files/X1 Pro/etc/uci-defaults/99_setup_wan
#!/bin/sh

uci set network.wan.metric=10

uci commit network.wan
