# files/X1 Pro/test/etc/uci-defaults/99_setup_firewall
#!/bin/sh

uci add_list firewall.@zone[1].network=mobiledata
uci set firewall.@zone[1].conntrack=1

uci commit firewall
