# files/T1/etc/init.d/startmodem
#!/bin/sh /etc/rc.common
# Copyright (C) 2020 Telco Antennas Pty Ltd

START=96

USE_PROCD=1

start_service() {
	ifup mobiledata
}

