# files/T1/etc/uci-defaults/9992_post_registration_hook
#!/bin/sh
# workaround to make default post-registration-hook executable
chmod +x /etc/openwisp/post-registration-hook
