# files/T1/etc/uci-defaults/999_setup_openwisp
#!/bin/sh
chmod +x /etc/openwisp/post-registration-hook