# files/X1 Pro/etc/crontabs/root
# Once a year, refresh the SSL certificates
0 0 1 */12 * /usr/bin/refresh_ssl_certificates.sh
