# files/X1 Pro/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file sets up your Telco Electronics device.
#
#
# These commands are required by your Telco Electronics device.
#
# BEGINNING of required commands.
#
sleep 20
/usr/sbin/irqbalance --oneshot
#/usr/bin/set_modem_in_10D_mode.sh
ifup mobiledata
#
# END of required commands.
#
#
exit 0
