# private_feeds/luci-theme-dark/root/etc/uci-defaults/30_luci-theme-dark
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.Dark=/luci-static/dark
		set luci.main.mediaurlbase=/luci-static/dark
		commit luci
	EOF
fi

exit 0
