# files/T1/files/etc/init.d/ledlistener
#!/bin/sh /etc/rc.common

START=99

start() {
	( /usr/bin/led_listener.sh )&
}

stop() {
	killall /usr/bin/led_listener.sh
}
